# systemd/graphical-idle.target
[Unit]
Description=Graphical session idle
Conflicts=graphical-unidle.target
