# systemd/graphical-unlock.target
[Unit]
Description=Graphical session unlock
Conflicts=graphical-lock.target
StopWhenUnneeded=yes
