# systemd/graphical-lock.target
[Unit]
Description=Graphical session lock
Conflicts=graphical-unlock.target
