# systemd/graphical-unidle.target
[Unit]
Description=Graphical session unidle
Conflicts=graphical-idle.target
StopWhenUnneeded=yes
