# systemd/sessiond-session.target
[Unit]
Description=sessiond session
Wants=sessiond.service
BindsTo=graphical-session.target
