# systemd/sessiond-session.target
[Unit]
Description=sessiond session
Requires=sessiond.service
BindsTo=graphical-session.target
