# systemd/user-sleep.target
[Unit]
Description=User sleep
StopWhenUnneeded=yes
