# systemd/user-sleep-finished.target
[Unit]
Description=User sleep finished
Conflicts=user-sleep.target
StopWhenUnneeded=yes
