# systemd/user-sleep.target
[Unit]
Description=User sleep
Conflicts=user-sleep-finished.target
