# systemd/user-shutdown.target
[Unit]
Description=User shutdown
StopWhenUnneeded=yes
